# supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=systemctl start httpd
